# server/openai-audio.service
[Unit]
Description=OpenAI Audio HTTP proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/xiaozhi
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/openai-audio
ExecStart=/opt/xiaozhi/.venv/bin/python /opt/xiaozhi/server/openai_audio_server.py --host 0.0.0.0 --port 8010
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
